# bdsd.service
[Unit]
Description=Bobaos Datapoint Sdk Daemon

[Service]
ExecStart=/usr/bin/env bdsd.sock

[Install]
WantedBy=default.target
